# setup/appupdates/crontab
00    09    *     *     *     appupdates-notify
